# src/docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid
childlogdir=/tmp
; 任意: conf.d を取り込む場合
; [include]
; files = /etc/supervisor/conf.d/*.conf

[program:socat]
command=/usr/bin/socat TCP-LISTEN:27124,bind=127.0.0.1,fork,reuseaddr TCP:host.docker.internal:27124
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1

[program:devshell]
command=/bin/zsh -l
user=node
environment=HOME="/home/node",USER="node"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/dev/fd/1
